# cookbooks/fuseesb/templates/default/fuse-esb_init.erb
#!/bin/bash
#
# Init script for Fuse ESB Enterprise
#
# Written by E Camden Fisher <camden.fisher@yale.edu>
#
# chkconfig: - 98 02
# description: Fuse Manamgement Console
#

source /etc/rc.d/init.d/functions

#########################################
# Sane defaults
MAXPERM="256M"
MAXHEAP="1024m"
MINHEAP="1024m"
# location of cacerts file
#CACERTS="/usr/local/etc/cacerts"
#########################################

# location of sysconfig file
[ -L $0 ] && SYSCONFIG="/etc/sysconfig/$(basename $(readlink $0))" || SYSCONFIG="/etc/sysconfig/$(basename $0)"

### Exit if sysconfig doesn't exist
if [ ! -r ${SYSCONFIG} ]; then
  echo "sysconfig: ${SYSCONFIG} not readable"
  exit 1
fi

### Read configuration
[ -r "$SYSCONFIG" ] && source "$SYSCONFIG"

# Set some karaf  variables
export KARAF_HOME=${KARAF_HOME}

# Set runtime Karaf options
export KARAF_OPTS="-Dkaraf.name=${KARAF_NAME} ${ADDLKRFOPTS}"

# Set some Java variables
export JAVA_HOME=${JAVA_HOME}
export JAVA_OPTS="-XX:MaxPermSize=${MAXPERM}  -Xms${MINHEAP} -Xmx${MAXHEAP} ${ADDLJVAOPTS}" 

# Set some normal env variables
export HOME=/home/${KARAF_USER}
export PATH=${JAVA_HOME}/bin:${PATH}

start(){
  /bin/echo "Starting Fuse ESB Enterprise ..."

  /usr/bin/sudo -b -E -u ${KARAF_USER} -- ${KARAF_HOME}/bin/start ${KARAF_OPTS} > /dev/null 2>&1

  sleep 10

}

stop(){
  /bin/echo "Stopping Fuse ESB Enterprise ..."

  /usr/bin/sudo -b -E -u ${KARAF_USER} -- ${KARAF_HOME}/bin/stop ${KARAF_OPTS} > /dev/null 2>&1

  /bin/sleep 10
  
}

restart(){
  stop

  /bin/sleep 30

  start
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  status)
        status
        ;;
  restart)
        restart
        ;;
  *)
        /bin/echo "Usage: $0 {start|stop|restart|status}"
        exit 1
esac

