# cookbooks/amq/templates/default/jboss-a-mq_init.erb
#!/bin/bash
#
# Init script for JBoss A-MQ
#
# Written by E Camden Fisher <camden.fisher@yale.edu>
#
# chkconfig: - 98 02
# description: JBoss A-MQ
#

source /etc/rc.d/init.d/functions

#########################################
# Sane defaults
MAXPERM="256M"
MAXHEAP="1024m"
MINHEAP="1024m"
#########################################

# location of sysconfig file
[ -L $0 ] && SYSCONFIG="/etc/sysconfig/$(basename $(readlink $0))" || SYSCONFIG="/etc/sysconfig/$(basename $0)"

### Exit if sysconfig doesn't exist
if [ ! -r ${SYSCONFIG} ]; then
  echo "sysconfig: ${SYSCONFIG} not readable"
  exit 1
fi

### Read configuration
[ -r "$SYSCONFIG" ] && source "$SYSCONFIG"

# Set some karaf  variables
export KARAF_HOME=${KARAF_HOME}

# Set some Java variables
export JAVA_HOME=${JAVA_HOME}
export JAVA_OPTS="-XX:MaxPermSize=${MAXPERM}  -Xms${MINHEAP} -Xmx${MAXHEAP} ${ADDLJVAOPTS}" 

# Set some normal env variables
export HOME=/home/${KARAF_USER}
export PATH=${PATH}:${JAVA_HOME}/bin

start(){
  FOO=`/usr/bin/pgrep -f -u ${KARAF_USER} ${KARAF_HOME}`

  if [ "$FOO" != "" ]
  then
    /bin/echo "Fuse ESB is already running..."
  else
    /bin/echo "Starting JBoss A-MQ ..."
    /usr/bin/sudo -b -E -u ${KARAF_USER} -- ${KARAF_HOME}/bin/start > /dev/null 2>&1
    sleep 20
  fi
}

stop(){
  FOO=`/usr/bin/pgrep -f -u ${KARAF_USER} ${KARAF_HOME}`
  if [ "$FOO" != "" ]
  then
    /bin/echo "Stopping JBoss A-MQ ..."
    /usr/bin/sudo -b -E -u ${KARAF_USER} -- ${KARAF_HOME}/bin/stop > /dev/null 2>&1
    ### wait a bit
    /bin/sleep 10
  else
    /bin/echo "Fuse ESB is not running..."
  fi
}

status(){

  #FOO=`/usr/bin/sudo -b -E -u ${KARAF_USER} -- ${KARAF_HOME}/bin/admin list | grep Started`
  FOO=`/usr/bin/pgrep -f -u ${KARAF_USER} ${KARAF_HOME}`

  if [ "$FOO" != "" ]
  then
    /bin/echo "jboss-a-mq is running"
    exit 0
  else
    /bin/echo "jboss-a-mq is stopped"
    exit 3
  fi

}

restart(){
  stop
  ### wait a bit
  /bin/sleep 30
  start
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  status)
        status
        ;;
  restart)
        restart
        ;;
  *)
        /bin/echo "Usage: $0 {start|stop|restart|status}"
        exit 1
esac

exit 0

