# cookbooks/jenkins/templates/default/jenkins_init.erb
#!/bin/bash
#
# Init script for Jenkins
#
# Written by E Camden Fisher <camden.fisher@yale.edu>
#
# chkconfig: - 98 02
# description: Jenkins CI
#

source /etc/rc.d/init.d/functions

# location of sysconfig file
[ -L $0 ] && SYSCONFIG="/etc/sysconfig/$(basename $(readlink $0))" || SYSCONFIG="/etc/sysconfig/$(basename $0)"

### Exit if sysconfig doesn't exist
if [ ! -r ${SYSCONFIG} ]; then
  echo "sysconfig: ${SYSCONFIG} not readable"
  exit 1
fi

### Read configuration
[ -r "$SYSCONFIG" ] && source "$SYSCONFIG"

export HOME=/home/${JENKINS_USER}
export JAVA_HOME=${JAVA_HOME}
export JENKINS_HOME=${JENKINS_HOME}
export JENKINS_BASE=${JENKINS_BASE}
export JENKINS_HTTPPORT=${JENKINS_HTTPPORT}
export JENKINS_AJPPPORT=${JENKINS_AJPPORT}
export PATH=${PATH}:${JAVA_HOME}/bin

start(){
  /bin/echo "Starting Jenkins..."

  ### fireup node
  /usr/bin/sudo -b -E -u ${JENKINS_USER} -- ${JAVA_HOME}/bin/java -jar ${JENKINS_BASE}/jenkins.war --httpPort=${JENKINS_HTTPPORT} --ajp13Port=${JENKINS_AJPPORT} >${JENKINS_HOME}/jenkins.log 2>&1 &
  sleep 10
}


stop(){
  /bin/echo "Stopping Jenkins ..."
  ### attempt a graceful pkill
  /usr/bin/pkill -f "jenkins.war" -u ${JENKINS_USER} 
  ### wait a bit
  /bin/sleep 15
  ### signal 9 any processes left-over
  /usr/bin/pkill -9 -f "jenkins.war" -u ${JENKINS_USER}

  /bin/sleep 5

}

restart(){
  stop
  ### wait a bit
  /bin/sleep 60
  start
}

case "$1" in
  start)
        start
        ;;
  start_debug)
        start_debug
        ;;
  stop)
        stop
        ;;
  restart)
        restart
        ;;
  *)
        /bin/echo "Usage: jenkins {start|start_debug|stop|restart}"
        exit 1
esac

exit 0

